# .config/systemd/user/session.target
[Unit]
Description="Session Target"
After=graphical-session-pre.target
WantedBy=graphical.target
Wants=drive.timer
Wants=fcitx5.service
